# with-supervisord/heartbeat.conf
[program:heartbeat]
command = /bin/bash /vagrant/with-supervisord/heartbeat.sh
autorestart = true
autostart = false
startsecs=10
startretries=3
