# with-systemd/heartbeat.service
[Unit]
Description=heartbeat example

[Service]
ExecStart=/bin/bash /vagrant/with-systemd/heartbeat.sh
Restart=always
